# config/crontab
# crontab for gitolite jobs

0 2 * * * /usr/local/bin/gitolite-mirror
