# config/crontab
# crontab for gitolite jobs

0 */2 * * *  su -c "/usr/local/bin/gitolite-mirror" $GIT_USER  > /var/log/gitolite-mirror.log
